# conf/etc/services.d/01-postfix/run
#!/bin/sh
NAME="POSTFIX"
COMMAND="postfix start-fg"
WAIT="/tmp/postfix"

exec service-handler "$NAME" "$COMMAND" "$WAIT"
