# conf/etc/services.d/02-dovecot/run
#!/bin/sh
NAME="DOVECOT"
COMMAND="dovecot -F"
WAIT="/tmp/dovecot"

exec service-handler "$NAME" "$COMMAND" "$WAIT"
